# l1price.service
[Unit]
Description=l1price - GenesisL1 (L1) price API
After=network-online.target
Wants=network-online.target

[Service]
User=ivan
ExecStart=/usr/bin/python3 /home/ivan/l1price.py
Environment=PORT=8585
Environment=CACHE_TTL=30
Restart=on-failure
RestartSec=10
# light hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
